# rootfs/etc/s6-overlay/s6-rc.d/02-piaware-log-prepare/up
if { mkdir -p /var/log/piaware }
if { chown nobody:nogroup /var/log/piaware }
chmod 02755 /var/log/piaware