# etc/services.d/dump1090-gain-stats/run
#!/usr/bin/with-contenv sh
sleep 3600
set -eo pipefail
/usr/local/bin/percentage_strong_messages.sh 2>&1 | awk '{print "[dump1090-gain-stats] " $0}'

